# led_server.service
[Unit]
Description=Start developer server

[Service]
ExecStart=/home/pi/shared/led_server/run_server.py
